# .config/systemd/user/nextcloud.service
[Unit]
Description=Nextcloud Synchronizers

[Service]
ExecStart=/home/shahzebs/.scripts/ncsync
Restart=always

[Install]
WantedBy=default.target
